# config/launchd/com.crows-nest.rss-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crows-nest.rss-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mriechers/Developer/second-brain/crows-nest/.venv/bin/python</string>
        <string>/Users/mriechers/Developer/second-brain/crows-nest/pipeline/rss_listener.py</string>
        <string>--refresh</string>
    </array>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>StandardOutPath</key>
    <string>/Users/mriechers/.local/share/crows-nest/logs/rss-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mriechers/.local/share/crows-nest/logs/rss-refresh.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
