# config/launchd/com.crows-nest.summarizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crows-nest.summarizer</string>

    <key>Comment</key>
    <string>Generate Obsidian notes from transcribed Crow's Nest content every 2 hours</string>

    <key>ProgramArguments</key>
    <array>
        <string>/PATH/TO/crows-nest/.venv/bin/python</string>
        <string>/PATH/TO/crows-nest/pipeline/summarizer.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/PATH/TO/crows-nest/pipeline</string>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>/PATH/TO/crows-nest/logs/launchd-crows-nest-summarizer.log</string>

    <key>StandardErrorPath</key>
    <string>/PATH/TO/crows-nest/logs/launchd-crows-nest-summarizer.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>

    <key>LowPriorityBackgroundIO</key>
    <true/>
</dict>
</plist>
