# config/launchd/com.crows-nest.processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crows-nest.processor</string>

    <key>Comment</key>
    <string>Process pending Crow's Nest links every 30 minutes — download media, transcribe, scrape articles</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mriechers/Developer/second-brain/crows-nest/.venv/bin/python</string>
        <string>/Users/mriechers/Developer/second-brain/crows-nest/pipeline/processor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mriechers/Developer/second-brain/crows-nest/pipeline</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/mriechers/.local/share/crows-nest/logs/processor.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mriechers/.local/share/crows-nest/logs/processor.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>USER</key>
        <string>mriechers</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>

    <key>LowPriorityBackgroundIO</key>
    <true/>
</dict>
</plist>
