# config/launchd/com.crows-nest.archiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crows-nest.archiver</string>

    <key>Comment</key>
    <string>Archive processed Crow's Nest content to Cloudflare R2 daily at 3 AM</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mriechers/Developer/second-brain/crows-nest/.venv/bin/python</string>
        <string>/Users/mriechers/Developer/second-brain/crows-nest/pipeline/archiver.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mriechers/Developer/second-brain/crows-nest/pipeline</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mriechers/.local/share/crows-nest/logs/archiver.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mriechers/.local/share/crows-nest/logs/archiver.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>USER</key>
        <string>mriechers</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>

    <key>LowPriorityBackgroundIO</key>
    <true/>
</dict>
</plist>
