# config/launchd/com.crows-nest.listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crows-nest.listener</string>

    <key>Comment</key>
    <string>Poll Signal for incoming URLs every 5 minutes and save to Crow's Nest pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/PATH/TO/crows-nest/.venv/bin/python</string>
        <string>/PATH/TO/crows-nest/pipeline/signal_listener.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/PATH/TO/crows-nest/pipeline</string>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/PATH/TO/crows-nest/logs/launchd-crows-nest-listener.log</string>

    <key>StandardErrorPath</key>
    <string>/PATH/TO/crows-nest/logs/launchd-crows-nest-listener.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
