# config/systemd/crows-nest-archiver.service
[Unit]
Description=Crow's Nest Archiver — archive processed content to Cloudflare R2
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/crows-nest/.venv/bin/python /opt/crows-nest/pipeline/archiver.py
WorkingDirectory=/opt/crows-nest/pipeline
Environment=CROWS_NEST_HOME=/opt/crows-nest
Environment=OBSIDIAN_VAULT=/data/obsidian/MarkBrain
Environment=MEDIA_ROOT=/data/media/crows-nest
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
